# contrib/umirr.socket
[Unit]
Description=umirr socket

[Socket]
ListenStream=/run/umirr/socket

[Install]
WantedBy=sockets.target
