# contrib/umirr@.socket
[Unit]
Description=umirr daemon instance socket

[Socket]
ListenStream=/run/umirr/%i-socket

[Install]
WantedBy=sockets.target
